# chat-server/chronograph-admin.service
[Unit]
Description=Chronograph Admin Panel (Tailscale only)
After=network.target chronograph-chat.service
Wants=chronograph-chat.service

[Service]
WorkingDirectory=/opt/chronograph-chat
User=redmond
ExecStart=/opt/chronograph-chat/venv/bin/python3 /opt/chronograph-chat/admin_panel.py
Restart=on-failure
RestartSec=5
Environment=CG_DB=/media/redmond/LargeHD/ChronoGraph-DB/chat.db

[Install]
WantedBy=default.target
